# build/rooted/etc/cron.d/letsencrypt-renew-nginx
01 02,14 * * * /etc/cron.daily/letsencrypt-renew-nginx
